# schedules/com.example.maestro-nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.example.maestro-nightly</string>

    <key>WorkingDirectory</key>
    <string>/path/to/project-root</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>npm run test:nightly</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>2</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
      <key>ENVIRONMENT</key>
      <string>nightly</string>
      <key>DEVICE_NAME</key>
      <string>Pixel_7_API_35</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/path/to/project-root/reports/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/project-root/reports/launchd.err.log</string>
  </dict>
</plist>
